# old/supervisord.conf
[supervisord]
nodaemon=false
 
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
 
 
[program:mysite]
command=python /mysite/manage.py runserver 0.0.0.0:80
user=root
autostart=true
directory=/usr/bin/
 
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
